# kcp.conf
[program:kcptun]
command=/opt/kcptun/server_linux_amd64 -l :29900 -t 127.0.0.1:8989 -crypt "none" --mtu 1400 --sndwnd 2048 --rcvwnd 2048 --mode "fast3"  --nocomp "true"
