# ssr.conf
[program:shadowsocksr]
priority = 100
command = python /ssr/shadowsocks/server.py  -c /opt/shadowsocks/user-config.json >/dev/nul
